# systemd/spiral-watcher.path
[Path]
PathModified=/home/spiral/bin/spiral
[Install]
WantedBy=multi-user.target